# templates/gnatsd-supervisor.conf
[program:gnatsd]
command={{ gnatsd_install_dir }}/gnatsd -c {{ gnatsd_conf_file }}
numprocs=1
autostart=true
autorestart=unexpected
